# supervisord.conf
[program:confluence]
command = /opt/atlassian/confluence/bin/start-confluence.sh 
autorestart = true